# php/crontab
0 */3 * * * /usr/local/bin/php /var/www/html/utils/parse.php
