# php/crontab
0 */3 * * * cd /var/www/html/utils && /usr/local/bin/php /var/www/html/utils/parse.php